# deploy/stockcalc.service
[Unit]
Description=StockCalc.Web ASP.NET Core service
After=network.target

[Service]
WorkingDirectory=/var/www/stockcalc
ExecStart=/usr/bin/dotnet /var/www/stockcalc/StockCalc.Web.dll
User=www-data
Group=www-data
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=ASPNETCORE_URLS=http://127.0.0.1:5001
Restart=always
RestartSec=5
KillSignal=SIGINT
SyslogIdentifier=stockcalc
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
